# com.lobster.office.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lobster.office</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alisa/.openclaw/workspace/projects/openclaw-office/start.sh</string>
        <string>-d</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/alisa/.openclaw/workspace/projects/openclaw-office</string>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/alisa/.openclaw/workspace/projects/openclaw-office/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/alisa/.openclaw/workspace/projects/openclaw-office/logs/launchd_stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/alisa</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LAUNCHD_MODE</key>
        <string>true</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>